# scripts/systemd/apotheon-connector-daily-brief.service
[Unit]
Description=Apotheon Website Connector - Generate Daily Brief
After=network.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/apotheon/connector.env
WorkingDirectory=/opt/dgx-bundle
ExecStart=/usr/bin/curl -sS -X POST http://127.0.0.1:8090/daily-brief \
  -H "Authorization: Bearer ${CONNECTOR_ADMIN_TOKEN}" \
  -H "Content-Type: application/json" \
  -d '{"goals": [], "audience": null}'
ExecStart=/bin/bash -lc "/opt/dgx-bundle/scripts/maintenance/connector-cleanup.sh"
