# scripts/systemd/apotheon-connector-daily-brief.timer
[Unit]
Description=Apotheon Website Connector - Daily Brief Timer

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
